# gemini-custom-mcp.service
[Unit]
Description=PeeperFrog Create MCP Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=peeperfrog
WorkingDirectory=/home/peeperfrog/peeperfrog-create
ExecStart=/usr/bin/python3 /home/peeperfrog/peeperfrog-create/peeperfrog-create-mcp/src/image_server.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Load environment from file (keeps secrets out of unit file)
EnvironmentFile=/home/peeperfrog/peeperfrog-create/.env

[Install]
WantedBy=multi-user.target
